# deploy/systemd/ollama.service
[Unit]
Description=Ollama LLM Server (local tier - qwen2.5:1.5b)
Documentation=https://ollama.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment="OLLAMA_HOST=0.0.0.0"
ExecStart=/usr/local/bin/ollama serve
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=default.target
