# deploy/systemd/tool-broker.service
[Unit]
Description=Smart Home Tool Broker (FastAPI)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Smart_Home
ExecStart=%h/Smart_Home/.venv/bin/uvicorn tool_broker.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
Environment="PATH=%h/Smart_Home/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
